# mms-agent/files/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="export">
  <service name="pkgsrc/mms-agent" type="service" version="0">
    <create_default_instance enabled="false"/>
    <single_instance/>
    <dependency name="fs-local" grouping="require_all" restart_on="none" type="service">
      <service_fmri value="svc:/system/filesystem/local"/>
    </dependency>
    <dependency name="network" grouping="require_all" restart_on="restart" type="service">
      <service_fmri value="svc:/milestone/network"/>
    </dependency>
    <dependency name="config-file" grouping="require_all" restart_on="refresh" type="path">
      <service_fmri value="file://localhost@PREFIX@/share/mms-agent/settings.py"/>
    </dependency>
    <exec_method name="start" type="method" exec="@PREFIX@/@SMF_METHOD_FILE@ start" timeout_seconds="60">
      <method_context working_directory="@MONGODB_DBPATH@">
        <method_credential user="@MONGODB_USER@" group="@MONGODB_GROUP@"/>
      </method_context>
    </exec_method>
    <exec_method name="stop" type="method" exec=":kill" timeout_seconds="60"/>
    <template>
      <common_name>
        <loctext xml:lang="C">MongoDB MMS Agent</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>

